# data/kagenti-generated/services/kiro-qms.service
[Unit]
Description=QMS monster shadow agent
After=network.target

[Service]
Type=simple
User=kiro
Group=kiro
SELinuxContext=system_u:system_r:svc_kiro_qms_t:s0

ExecStart=/mnt/data1/zones/42/uucp/kiro-qms/bin/qms
Restart=always
RestartSec=10
Environment="MONSTER_SHARD=19"
Environment="IPV6_ADDR=fd00::71:19:0"

# Security (monster-zone gold standard)
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only

# Network isolation
IPAddressAllow=10.71.0.0/16
IPAddressAllow=fd00::71:0:0/48
IPAddressAllow=127.0.0.0/8
IPAddressAllow=::1/128

# Resource limits
MemoryMax=2G
TasksMax=64

[Install]
WantedBy=multi-user.target