# data/kagenti-generated/services/moltis.service
[Unit]
Description=Moltbot distributed agent (zone66)
After=network.target

[Service]
Type=simple
User=moltis
Group=moltis
SELinuxContext=system_u:system_r:svc_moltis_t:s0

ExecStart=/opt/moltis/bin/moltis-agent
Restart=always
RestartSec=10
Environment="PORT=8130"
Environment="MONSTER_SHARD=65"
Environment="IPV6_ADDR=fd00::71:65:0"

# Security (monster-zone gold standard)
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only

# Network isolation
IPAddressAllow=10.71.0.0/16
IPAddressAllow=fd00::71:0:0/48
IPAddressAllow=127.0.0.0/8
IPAddressAllow=::1/128

# Resource limits
MemoryMax=2G
TasksMax=64

[Install]
WantedBy=multi-user.target