# data/kagenti-generated/services/pastebin.service
[Unit]
Description=Stego pastebin with WASM decoder
After=network.target

[Service]
Type=simple
User=mdupont
Group=mdupont
SELinuxContext=system_u:system_r:svc_pastebin_t:s0

ExecStart=/opt/pastebin/bin/pastebin-server
Restart=always
RestartSec=10
Environment="PORT=8150"
Environment="MONSTER_SHARD=38"
Environment="IPV6_ADDR=fd00::71:38:0"

# Security (monster-zone gold standard)
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only

# Network isolation
IPAddressAllow=10.71.0.0/16
IPAddressAllow=fd00::71:0:0/48
IPAddressAllow=127.0.0.0/8
IPAddressAllow=::1/128

# Resource limits
MemoryMax=2G
TasksMax=64

[Install]
WantedBy=multi-user.target