# data/kagenti-generated/services/wg-stego-tunnel.service
[Unit]
Description=WireGuard tunnel hidden in cat memes
After=network.target

[Service]
Type=simple
User=wgtunnel
Group=wgtunnel
SELinuxContext=system_u:system_r:svc_wg_stego_tunnel_t:s0

ExecStart=/opt/wg-stego/bin/wg-stego-tunnel
Restart=always
RestartSec=10
Environment="PORT=8160"
Environment="MONSTER_SHARD=7"
Environment="IPV6_ADDR=fd00::71:7:0"

# Security (monster-zone gold standard)
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only

# Network isolation
IPAddressAllow=10.71.0.0/16
IPAddressAllow=fd00::71:0:0/48
IPAddressAllow=127.0.0.0/8
IPAddressAllow=::1/128

# Resource limits
MemoryMax=2G
TasksMax=64

[Install]
WantedBy=multi-user.target