# data/kagenti-generated/services/zos-server.service
[Unit]
Description=zkWASM executor + plugin hub
After=network.target

[Service]
Type=simple
User=zos
Group=zos
SELinuxContext=system_u:system_r:svc_zos_server_t:s0

ExecStart=/opt/zos/bin/zos-server
Restart=always
RestartSec=10
Environment="PORT=8190"
Environment="MONSTER_SHARD=34"
Environment="IPV6_ADDR=fd00::71:34:0"

# Security (monster-zone gold standard)
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only

# Network isolation
IPAddressAllow=10.71.0.0/16
IPAddressAllow=fd00::71:0:0/48
IPAddressAllow=127.0.0.0/8
IPAddressAllow=::1/128

# Resource limits
MemoryMax=2G
TasksMax=64

[Install]
WantedBy=multi-user.target